# mexc-research-collector.service
[Unit]
Description=MEXC Research Collector Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/mexc-research-collector
EnvironmentFile=/opt/mexc-research-collector/.env
ExecStart=/opt/mexc-research-collector/venv/bin/python /opt/mexc-research-collector/run.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
