# mexc-research-collector.service
[Unit]
Description=ChatGPT Scan Bot 30d
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/chatgpt-scan-bot
EnvironmentFile=/opt/chatgpt-scan-bot/.env
ExecStart=/opt/chatgpt-scan-bot/venv/bin/python /opt/chatgpt-scan-bot/run.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
